# vaio/vaio1.timer
[Unit]
Description=Run webscraper daily at 6 am

[Timer]
OnCalendar=*-*-* 06:05:00
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target

